# scripts/super-price-il.timer
[Unit]
Description=super-price-il price ingest — weekly cadence with on-demand catch-up

[Timer]
# Hourly tick. The script itself uses --skip-recent 7, so a chain scraped
# within 7 days is a no-op (very fast). Anything older / never-scraped
# runs immediately on the next tick.
OnCalendar=*-*-* *:15:00
Persistent=true
Unit=super-price-il.service

[Install]
WantedBy=timers.target
